# systemd/psk-recorder@.service
[Unit]
Description=psk-recorder FT4/FT8 spot recorder for radiod %i
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=notify
User=pskrec
Group=pskrec

WatchdogSec=120
NotifyAccess=main

EnvironmentFile=-/etc/sigmond/coordination.env
EnvironmentFile=-/etc/psk-recorder/env/%i.env

WorkingDirectory=/opt/psk-recorder

ExecStartPre=+/bin/mkdir -p /var/lib/psk-recorder/%i/ft8
ExecStartPre=+/bin/mkdir -p /var/lib/psk-recorder/%i/ft4
ExecStartPre=+/bin/mkdir -p /var/log/psk-recorder
ExecStartPre=+/usr/bin/chown -R pskrec:pskrec /var/lib/psk-recorder/%i
ExecStartPre=+/usr/bin/chown -R pskrec:pskrec /var/log/psk-recorder
# /var/lib/hs-uploader is SHARED across HamSCI client users (pskrec,
# wsprdaemon, ...) — its setgid+group-write shape is owned by
# tmpfiles.d/hs-uploader.conf (root:sigmond, mode 02775).  Earlier
# versions of this unit had:
#
#   ExecStartPre=+/bin/mkdir -p /var/lib/hs-uploader
#   ExecStartPre=+/usr/bin/chown pskrec:pskrec /var/lib/hs-uploader
#
# Both wrong: the first redundant with tmpfiles.d, the second hostile —
# it re-stomped the dir's group from sigmond to pskrec on every
# psk-recorder restart, locking wsprdaemon-client's wd-upload-hs out
# of the shared watermarks.db with "attempt to write a readonly
# database".  Observed on bee1 2026-05-14.  Removed.

ExecStart=/opt/psk-recorder/venv/bin/python3 -m psk_recorder.cli daemon \
    --config /etc/psk-recorder/psk-recorder-config.toml \
    --radiod-id %i

# 180s covers a typical 19-channel (FT8 × 10 bands + FT4 × 9 bands) startup
# plus headroom: each channel triggers a setpolicy round-trip with radiod,
# the initial multicast discovery sleeps 1s, and a fresh start after a
# stop adds ~10s of teardown.  60s was racy on a busy host.
TimeoutStartSec=180

Restart=always
RestartSec=5

# Logs go to the systemd journal (sigmond standard) — query with
# `journalctl -u psk-recorder@%i` or `smd log psk-recorder`.
StandardOutput=journal
StandardError=inherit
SyslogIdentifier=psk-recorder@%i

Nice=5
# 19+ concurrent jt9 children at ~60 MB RSS each saturates 512M.  jt9
# is a Fortran binary that mmaps its own working buffers per slot;
# this scales linearly with channel count and decoding depth.  2 GB
# leaves comfortable headroom for the standard FT8 × 10 + FT4 × 9
# layout plus the psk-recorder Python process itself.
MemoryMax=2G
MemorySwapMax=0

NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/psk-recorder /var/log/psk-recorder /var/lib/hs-uploader
ReadOnlyPaths=/etc/psk-recorder
PrivateNetwork=false

CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
